# com.trighub.webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trighub.webhook</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/trinhhao/trighub-webhook.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/trinhhao/trighub-webhook.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/trinhhao/trighub-webhook-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/trinhhao</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/trinhhao</string>
        <key>TRIGHUB_SECRET</key>
        <string></string>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
